# deploy/survey-bot-reminder.service.j2
[Unit]
Description=SerwisWin Survey Bot Reminder
After=network.target

[Service]
Type=oneshot
User=surveybot
Group=surveybot
Environment=DATA_DIR={{ data_dir }}
WorkingDirectory={{ bot_dir }}
ExecStart=/usr/bin/python3 {{ bot_dir }}/reminder.py
StandardOutput=journal
StandardError=journal
